# EmbeddedLinux/08.Init_Process/Services/myservice.service
[Unit]
Want=mynewtwork.network
[service]
Type=simple
ExecStart=ping -c 4 google.com
Restart=always
RestartSec=3
WorkDir=/usr/bin
